# emergent/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:postgres]
command=docker-entrypoint.sh postgres
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres_err.log
environment=POSTGRES_DB="%(ENV_POSTGRES_DB)s",POSTGRES_USER="%(ENV_POSTGRES_USER)s",POSTGRES_PASSWORD="%(ENV_POSTGRES_PASSWORD)s"

[program:redis]
command=redis-server --requirepass %(ENV_REDIS_PASSWORD)s --bind 127.0.0.1
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log

[program:backend]
command=node dist/main.js
directory=/app/backend
autostart=true
autorestart=true
priority=10
startretries=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
environment=NODE_ENV="production",PORT="8001"

[program:frontend]
command=node .next/standalone/server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=11
startretries=5
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend_err.log
environment=PORT="3000",HOSTNAME="0.0.0.0"
